# deploy/cc-prospects.service
[Unit]
Description=Common Crawl prospect collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/cc-links-scoring
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/cc-links-scoring/.venv/bin/python /opt/cc-links-scoring/multi_crawl.py --target-total 100000 --max-crawls 12 --state-dir /var/lib/cc-prospects/crawls --categories-file /opt/cc-links-scoring/categories.json --per-category-limit 5000 --max-parts 300 --max-per-domain 10 --min-score 50 --db /var/lib/cc-prospects/prospects.db --workers 64 --source s3 --progress-interval 60
Restart=on-failure
RestartSec=30
TimeoutStopSec=300
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
